# cauldron/supervisor_cauldron.conf
[program:cauldron]
command = /code/start_gunicorn.sh
user = root
stdout_logfile = /logs/gunicorn_supervisor.log
redirect_stderr = true